# run/visiobas_gateway.service
[Unit]
Description=VisioBAS Gateway

[Service]
User=root
Restart=on-failure
WorkingDirectory=/opt/visiobas-gateway
ExecStart=poetry shell && python3 /opt/visiobas-gateway/visiobas_gateway
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
Alias=visiobas_gateway.service